# deploy/crypto-lake.service
[Unit]
Description=Crypto Lake Data Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=Eschaton
WorkingDirectory=/home/Eschaton/crypto-lake
Environment="PATH=/home/Eschaton/crypto-lake/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/etc/default/crypto-lake
ExecStart=/home/Eschaton/crypto-lake/venv/bin/python main.py --mode orchestrate
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
